# ecnu_login.service.template
[Unit]
Description=login into ecnu network

[Service]
Type=oneshot
WorkingDirectory=$PWD
ExecStart=sh ./ecnu_login.sh
User=${USER}
StandardOutput=journal
StandardError=inherit
Group=systemd-journal
